# bash/.config/cx.darryl.bash/profile.d/10-bash.bash
file=${HOMEBREW_PREFIX:-/opt/homebrew}/etc/profile.d/bash_completion.sh

[[ -r $file ]] && source $file

unset file
