# iterm2/.config/cx.darryl.bash/profile.d/10-iterm2.bash
file=${XDG_DATA_HOME:-$HOME/.local/share}/iterm2/shell_integration.bash

[[ -e $file ]] && source $file

unset file
