# modules/.config/cx.darryl.bash/profile.d/01-modules.bash
source ${XDG_DATA_HOME:-$HOME/.local/share}/cx.darryl.modules/init/bash
