# less/.config/cx.darryl.bash/profile.d/10-less.bash
export LESSHISTFILE=${XDG_STATE_HOME:=$HOME/.local/state}/less/history
