# readline/.config/cx.darryl.bash/profile.d/00-readline.bash
export INPUTRC=${XDG_CONFIG_HOME:-$HOME/.config}/readline/inputrc
