# amazon-q/.config/cx.darryl.bash/profile.d/00-amazon-q.bash
source $HOME/Library/Application\ Support/amazon-q/shell/bash_profile.pre.bash
