# amazon-q/.config/cx.darryl.bash/profile.d/99-amazon-q.bash
source $HOME/Library/Application\ Support/amazon-q/shell/bash_profile.post.bash
